# docker/etc/s6-overlay/s6-rc.d/svc-sshd/run
#!/usr/bin/with-contenv bash

set -e

exec /usr/sbin/sshd -D
