# docker/root/etc/s6-overlay/s6-rc.d/svc-flaskfarm/run
#!/usr/bin/with-contenv bash

set -e

source ${FF_DATA}/svc
